# AikoR-linux-64/DEBIAN/prerm
#!/bin/bash

unlink /etc/AikoR/aiko.yml
systemctl stop AikoR
systemctl disable AikoR
